# scripts/edge/linux/python/edge-monitoring.service
[Unit]
Description=Edge Monitoring Daemon
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/monitoring
ExecStart=/usr/bin/python3 /opt/monitoring/edge_daemon.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
